# templates/launchagents/com.esko.llama-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.esko.llama-server</string>
	<key>Disabled</key>
	<true/>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>__LLAMA_SERVER_BIN__</string>
		<string>--model</string>
		<string>__MODEL_PATH__</string>
		<string>--host</string>
		<string>__LLAMA_HOST__</string>
		<string>--port</string>
		<string>__LLAMA_PORT__</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>__HOME__</string>
	</dict>
	<key>StandardOutPath</key>
	<string>__LLAMA_LOG_DIR__/llama-server.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>__LLAMA_LOG_DIR__/llama-server.stderr.log</string>
</dict>
</plist>
